# deploy/systemd/housing-worker.service
[Unit]
Description=Housing portal — Celery worker
After=network.target redis-server.service housing-backend.service
Wants=redis-server.service

[Service]
Type=simple
User=housing
Group=housing
WorkingDirectory=/opt/housing/backend
EnvironmentFile=/opt/housing/backend/.env
ExecStart=/opt/housing/backend/.venv/bin/celery -A celery_worker.celery worker --loglevel=info
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
