# deploy/systemd/housing-beat.service
[Unit]
Description=Housing portal — Celery beat scheduler
After=network.target redis-server.service housing-worker.service
Wants=redis-server.service

[Service]
Type=simple
User=housing
Group=housing
WorkingDirectory=/opt/housing/backend
EnvironmentFile=/opt/housing/backend/.env
ExecStart=/opt/housing/backend/.venv/bin/celery -A celery_worker.celery beat --loglevel=info --schedule=/var/lib/housing/celerybeat-schedule
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
# beat needs to write its schedule file under /var/lib/housing.
ReadWritePaths=/var/lib/housing
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
